# templates/default/burpui.conf.erb
description "BURP User Interface"
start on net-device-up
stop on runlevel [06]
respawn
exec /usr/local/bin/burp-ui
